# calibre-kindle.cron
# Autogenerated crontab schedule for Calibre Kindle Server

SHELL=/bin/bash
PATH=/bin:/usr/bin:/usr/local/bin

# m h  d m y	command
0 6 * * *    /root/calibre-kindle-server/deliver.py
0 13 * * *    /root/calibre-kindle-server/deliver.py
0 20 * * *    /root/calibre-kindle-server/deliver.py
0,12,24,36,48 * * * *     /root/calibre-kindle-server/proces-recipes.py
