# calibre-kindle.cron
# Crontab schedule for Calibre Kindle Server

SHELL=/bin/bash
PATH=/bin:/usr/bin:/usr/local/bin

# m  h	d m y	command
  30 04	* * *	cd /root/calibre-kindle-server/ && deliver.py bbcnews
  00 13	* * *	cd /root/calibre-kindle-server/ && deliver.py bbcnews
  00 22	* * *	cd /root/calibre-kindle-server/ && deliver.py bbcnews

  00 06	* * *	cd /root/calibre-kindle-server/ && deliver.py acmtechnews
